# pkg/ember-notify.install
post_install() {
  cat <<'EOF'
==> Ember installed.
==> Enable the user service with:
    systemctl --user enable --now ember.service
==> Optional user config:
    cp /usr/share/ember/default.toml ~/.config/ember/config.toml
EOF
}

post_upgrade() {
  post_install
}
